# deployment/com.decifer.handoff-publisher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.decifer.handoff-publisher.plist
  Classification: documentation only / deployment template

  Intraday manifest refresh daemon for the Intelligence-First handoff.

  Runs handoff_publisher.py --mode controlled_activation every 10 minutes.
  This keeps data/live/current_manifest.json within its 15-minute expiry window,
  ensuring the bot can consume the prepared universe on every scan cycle.

  Without this daemon (or equivalent), the manifest expires 15 minutes after the
  last manual publish and bot_trading.py fails closed (universe = [], Track A skipped).

  Schedule: Every 600 seconds (10 minutes), all days.
  On load: fires immediately (RunAtLoad = true) so the manifest is fresh at startup.
  No broker calls, no LLM calls, no live data. Pure file read/write.

  INSTALLATION:
    cp deployment/com.decifer.handoff-publisher.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.decifer.handoff-publisher.plist

  EVIDENCE:
    Heartbeat: data/heartbeats/handoff_publisher.json
    Stdout:    /tmp/decifer-handoff-publisher.log
    Stderr:    /tmp/decifer-handoff-publisher.err
    Manifest:  data/live/current_manifest.json (expires_at refreshed each run)

  VERIFY:
    launchctl list | grep decifer.handoff-publisher
    python3 -c "import json; m=json.load(open('data/live/current_manifest.json')); print(m['expires_at'], m['handoff_enabled'])"

  DEPENDENCY: Reads from data/universe_builder/active_opportunity_universe_shadow.json.
  That file is updated daily by com.decifer.intelligence-pipeline (07:30).
  Between daily refreshes this daemon re-publishes the same universe with a fresh
  timestamp — correct behaviour, as universe content is stable intraday.

  FAILURE MODE: If the shadow universe is missing or invalid, the publisher fails
  closed and writes a .fail_{timestamp}.json diagnostic to data/live/. The bot
  continues running existing positions via Track B and forced-exit logic.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.decifer.handoff-publisher</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/python3.11</string>
    <string>handoff_publisher.py</string>
    <string>--mode</string>
    <string>controlled_activation</string>
  </array>

  <!-- Adjust this path to your local repo root -->
  <key>WorkingDirectory</key>
  <string>/Users/amitchopra/Desktop/decifer trading</string>

  <!-- Every 10 minutes, all days. -->
  <key>StartInterval</key>
  <integer>600</integer>

  <key>StandardOutPath</key>
  <string>/tmp/decifer-handoff-publisher.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/decifer-handoff-publisher.err</string>

  <key>KeepAlive</key>
  <false/>

  <!-- Fire immediately on load so the manifest is fresh at bot startup. -->
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
