# deployment/com.decifer.intelligence-pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.decifer.intelligence-pipeline.plist
  Classification: documentation only / deployment template

  Daily pre-market Intelligence-First pipeline refresh.

  Runs the full 3-step chain in sequence:
    1. run_intelligence_pipeline.py  — resolves candidates, economic state, theme activation
    2. universe_builder.py           — assembles shadow universe from fresh intelligence
    3. handoff_publisher.py          — publishes manifest + active universe to data/live/

  Schedule: Monday–Friday at 07:30 local time (pre-NYSE open).
  Runs independently of bot.py. No broker calls, no LLM calls, no live data.

  INSTALLATION:
    cp deployment/com.decifer.intelligence-pipeline.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.decifer.intelligence-pipeline.plist

  EVIDENCE:
    Stdout log: /tmp/decifer-intelligence-pipeline.log
    Stderr log: /tmp/decifer-intelligence-pipeline.err
    Each step also writes its own heartbeat/output:
      data/intelligence/economic_candidate_feed.json
      data/intelligence/daily_economic_state.json
      data/intelligence/theme_activation.json
      data/universe_builder/active_opportunity_universe_shadow.json
      data/live/current_manifest.json
      data/heartbeats/handoff_publisher.json

  VERIFY:
    launchctl list | grep decifer.intelligence-pipeline
    cat /tmp/decifer-intelligence-pipeline.log

  NOTE: On a fresh machine, run the chain once manually before installing:
    cd "decifer trading"
    python3 run_intelligence_pipeline.py
    python3 universe_builder.py
    python3 handoff_publisher.py --mode controlled_activation
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.decifer.intelligence-pipeline</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-c</string>
    <string>/usr/local/bin/python3.11 run_intelligence_pipeline.py &amp;&amp; /usr/local/bin/python3.11 universe_builder.py &amp;&amp; /usr/local/bin/python3.11 handoff_publisher.py --mode controlled_activation</string>
  </array>

  <!-- Adjust this path to your local repo root -->
  <key>WorkingDirectory</key>
  <string>/Users/amitchopra/Desktop/decifer trading</string>

  <!-- 07:30 Monday–Friday (pre-market) -->
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Weekday</key><integer>1</integer>
      <key>Hour</key><integer>7</integer>
      <key>Minute</key><integer>30</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>2</integer>
      <key>Hour</key><integer>7</integer>
      <key>Minute</key><integer>30</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>3</integer>
      <key>Hour</key><integer>7</integer>
      <key>Minute</key><integer>30</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>4</integer>
      <key>Hour</key><integer>7</integer>
      <key>Minute</key><integer>30</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>5</integer>
      <key>Hour</key><integer>7</integer>
      <key>Minute</key><integer>30</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>/tmp/decifer-intelligence-pipeline.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/decifer-intelligence-pipeline.err</string>

  <key>KeepAlive</key>
  <false/>

  <!-- Do not run at load — only fire on the scheduled calendar interval. -->
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
